# scripts/raspberry-pi/accuweather-culver-city.service
[Unit]
Description=AccuWeather Culver City Forecast Server
Documentation=https://github.com/your-repo/accuweather-culver-city
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/accuweather-culver-city
Environment=NODE_ENV=production
EnvironmentFile=-/home/pi/accuweather-culver-city/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
